# deploy/systemd/ounje-recipe-ingestion-worker.service
[Unit]
Description=Ounje private recipe ingestion worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ounje
WorkingDirectory=/opt/ounje
Environment=NODE_ENV=production
Environment=OUNJE_RUNTIME_ENV=production
Environment=OUNJE_ENABLE_AI_CALL_LOGGING=true
Environment=RECIPE_INGESTION_WORKER_CONCURRENCY=2
Environment=RECIPE_INGESTION_WORKER_WAKE_MODE=redis
Environment=RECIPE_INGESTION_WORKER_WAKE_TIMEOUT_MS=900000
Environment=YOUTUBE_DL_BINARY=/usr/local/bin/yt-dlp
EnvironmentFile=-/etc/ounje/ounje.env
ExecStart=/usr/bin/env node /opt/ounje/server/scripts/recipe_ingestion_worker.mjs --worker-id vm_recipe_ingest --batch-size 2 --wake-mode redis --wake-timeout-ms 900000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
