# deploy/systemd/ounje-automation-worker.service
[Unit]
Description=Ounje private automation worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ounje
WorkingDirectory=/opt/ounje
Environment=NODE_ENV=production
Environment=OUNJE_RUNTIME_ENV=production
EnvironmentFile=-/etc/ounje/ounje.env
ExecStart=/usr/bin/env node /opt/ounje/server/scripts/automation_worker.mjs --worker-id vm_automation_worker --lock-seconds 900
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
